# car_app/supervisord.conf
[supervisord]
nodaemon=true

[program:telegraf]
command=/usr/bin/telegraf
